# layers/meta-infotainment/recipes-graphics/wayland-layer-manager/files/wayland-layer-manager.service
[Unit]
Description=Wayland Layer Manager Compositor
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/wayland-layer-manager
Environment="XDG_RUNTIME_DIR=/run/user/0"
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
