# layers/meta-infotainment/recipes-graphics/ivi-compositor/files/ivi-compositor.service
[Unit]
Description=Qt Wayland Compositor Service
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/ivi-compositor
Environment="XDG_RUNTIME_DIR=/run/user/0"
Environment="QT_QPA_PLATFORM=eglfs"
Environment="WAYLAND_DEBUG=1"
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
