# layers/meta-infotainment/recipes-connectivity/can-setup/files/can-setup.service
[Unit]
Description=CAN interface setup
After=network.target

[Service]
Type=oneshot
ExecStart=/sbin/ip link set can0 up type can bitrate 500000 fd off
ExecStop=/sbin/ip link set can0 down
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
